# rootfs2/etc/inittab




























null::sysinit:/bin/mount -a
null::sysinit:/bin/hostname -F /etc/hostname
null::respawn:/bin/node /projeto/server.js
null::respawn:/bin/sh /led_off.sh
null::respawn:/bin/cttyhack /bin/login root
null::restart:/sbin/reboot

